# contrib/scripts/vma.service.in
[Unit]
Description=VMA Daemon. Version: @VERSION@-@VMA_LIBRARY_RELEASE@
After=network.target syslog.target
Requires=network.target

[Service]
Type=forking
Restart=on-failure
ExecStart=@prefix@/sbin/vma start
ExecStop=@prefix@/sbin/vma stop
ExecReload=@prefix@/sbin/vma restart
RestartForceExitStatus=1 SIGTERM

[Install]
WantedBy=multi-user.target
