# contrib/scripts/vma.service.in
[Unit]
Description=VMA Daemon
After=network.target

[Service]
Type=forking
Restart=on-failure
ExecStart=@prefix@/sbin/vmad
RestartForceExitStatus=1 SIGTERM

[Install]
WantedBy=multi-user.target
