# debian/postinst
#!/bin/bash

/sbin/ldconfig

if command -v systemctl >/dev/null 2>&1; then
    systemctl --no-reload enable vma.service >/dev/null 2>&1 || true
fi
